# src/ApiServer/supervisord.conf
[supervisord]
nodaemon=true

[program:dotnet]
command=dotnet ApiServer.dll

[program:sshd]
command=/usr/sbin/sshd -D
